# systemd/deepseek-cursor-pending-watcher.path
[Unit]
Description=Watch for pending Cursor DeepSeek URL updates
Documentation=https://github.com/cheapjourney/deepseek-cursor-setup

[Path]
PathExists=%h/.cache/deepseek-cursor-proxy/pending-base-url.txt

[Install]
WantedBy=default.target
