# systemd/update-cursor-deepseek-url.timer
[Unit]
Description=Retry Cursor DeepSeek URL update after Cloudflare Quick Tunnel start

[Timer]
OnBootSec=30s
OnUnitInactiveSec=60s
OnCalendar=minutely
AccuracySec=5s
Unit=update-cursor-deepseek-url.service
Persistent=true

[Install]
WantedBy=timers.target
