# systemd/update-cursor-deepseek-url.timer
[Unit]
Description=Retry Cursor DeepSeek URL update after Cloudflare Quick Tunnel start

[Timer]
OnBootSec=90s
OnUnitActiveSec=90s
AccuracySec=10s
Unit=update-cursor-deepseek-url.service
Persistent=true

[Install]
WantedBy=timers.target
